# contrib/GUI/vm/ubuntu22/secondsight.service
[Unit]
Description=Second Sight

[Service]
Type=simple
ExecStart=/usr/bin/java -jar /app/secondsight.jar

[Install]
WantedBy=multi-user.target
